# misc-minecraft/mc.cron
0 5 * * * /bin/bash /home/<USER_HERE>/backupmc.sh

